# rdumpfs-rota-daily.cron
#!/bin/bash
/usr/bin/rdumpfs-rota daily
