# rdumpfs-rota-monthly.cron
#!/bin/bash
/usr/bin/rdumpfs-rota rdumpfs monthly
/usr/bin/rdumpfs-rota ldumpfs monthly
