# rdumpfs-rota-weekly.cron
#!/bin/bash
/usr/bin/rdumpfs-rota weekly
