# rdumpfs-rota-monthly.cron
#!/bin/bash
/usr/bin/rdumpfs-rota monthly
