# etc/rdumpfs-rota-weekly.cron
#!/bin/bash
/usr/bin/rdumpfs-rota rdumpfs weekly
/usr/bin/rdumpfs-rota ldumpfs weekly
