# etc/rdumpfs-rota-daily.cron
#!/bin/bash
/usr/bin/rdumpfs-rota rdumpfs daily
/usr/bin/rdumpfs-rota ldumpfs daily
